# supervisor/arch/dbus.ini
[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
